# filepuller.service
[Service]
ExecStart=/opt/rheapdumps/filepuller.daemon
User=rheapdumps

[Install]
WantedBy=multi-user.target
