# remotefilewatch@.service
[Service]
ExecStart=/opt/rheapdumps/remotefilewatch.daemon %i
User=rheapdumps

[Install]
WantedBy=multi-user.target
